# ser2net.service.in
[Unit]
Description=Serial to network proxy
Requires=network.target
Documentation=man:ser2net(8)

[Service]
Type=simple
User=ser2net
Group=dialout
ExecStart=@ser2netbindir@/ser2net -n
ExecReload=kill -HUP $MAINPID
StandardOutput=null
Restart=on-failure

[Install]
WantedBy=multi-user.target
